# scripts/vilesql.service
[Unit]
Description=VileSQL Database Management Service
Documentation=https://github.com/imrany/vilesql
After=network.target network-online.target
Wants=network-online.target
RequiresMountsFor=/var/lib/vilesql

[Service]
Type=simple
User=vilesql
Group=vilesql
WorkingDirectory=/var/lib/vilesql

# Load environment variables before executing VileSQL
EnvironmentFile=/var/lib/vilesql/.env
ExecStart=/bin/bash -c 'source /var/lib/vilesql/.env && exec /usr/bin/vilesql --host=$HOST --port=5000'
ExecReload=/bin/kill -USR2 $MAINPID

Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=process

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
SystemCallArchitectures=native

# Allow VileSQL write access to its working directory
ReadWritePaths=/var/lib/vilesql

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=64M

# Logging (redirect to file)
StandardOutput=append:/var/log/vilesql.log
StandardError=append:/var/log/vilesql.log

[Install]
WantedBy=multi-user.target